# examples/supervisord/supervisord.example.conf
; examples/supervisord/supervisord.example.conf
; Reference supervisord layout for an opencode-mem deployment.
; Programs in order:
;   1. config-manager  — writes opencode-mem.jsonc, emits a readiness flag, exits
;   2. migrate         — runs schema migration once; exits with code 0
;   3. opencode        — main agent; depends on the readiness flag being present
;   4. qdrant (opt)    — sidecar for users pairing PostgresRecordStore + USearch (or future Qdrant)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:config-manager]
command=/usr/local/bin/config-manager.sh
autostart=true
autorestart=false
startretries=0
exitcodes=0
priority=10
stdout_logfile=/var/log/supervisor/config-manager.log
stderr_logfile=/var/log/supervisor/config-manager.err

[program:migrate]
command=npx --yes opencode-mem-migrate --to postgres --url %(ENV_DATABASE_URL)s
autostart=true
autorestart=false
startretries=0
exitcodes=0
priority=20
stdout_logfile=/var/log/supervisor/migrate.log
stderr_logfile=/var/log/supervisor/migrate.err

[program:opencode]
command=/usr/local/bin/wait-for-ready-then-run-opencode.sh
autostart=true
autorestart=true
startretries=10
priority=30
stdout_logfile=/var/log/supervisor/opencode.log
stderr_logfile=/var/log/supervisor/opencode.err

[program:qdrant]
command=/usr/local/bin/qdrant
autostart=false ; flip to true if you pair Postgres + USearch and want a future Qdrant backend
autorestart=true
priority=25
stdout_logfile=/var/log/supervisor/qdrant.log
stderr_logfile=/var/log/supervisor/qdrant.err
